# contrib/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable dcvix-agent.service 
systemctl start dcvix-agent.service 