# contrib/deb/prerm
#!/bin/sh
set -e

systemctl stop dcvix-agent.service || true
systemctl disable dcvix-agent.service || true 